# deploy/shitposter.timer
[Unit]
Description=Run shitposter daily

[Timer]
OnCalendar=RUN_SCHEDULE RUN_TIMEZONE
Persistent=true

[Install]
WantedBy=timers.target
